# deploy/gunicorn.service
# Systemd servis za Gunicorn
# Kopirajte u /etc/systemd/system/raspored.service
# Zatim: systemctl daemon-reload && systemctl enable raspored && systemctl start raspored

[Unit]
Description=FOOZOS Raspored Gunicorn
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/raspored
Environment="FLASK_ENV=production"
ExecStart=/var/www/raspored/venv/bin/gunicorn \
    --workers 3 \
    --timeout 120 \
    --access-logfile /var/log/raspored/access.log \
    --error-logfile /var/log/raspored/error.log \
    --bind 127.0.0.1:8000 \
    run:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
